# crontab
0 7 * * * ~/scripts/on.sh \n0 16 * * * ~/scripts/off.sh -e
